# llama.cpp/llama-cpp.service
[Unit]
Description=Llama.cpp OpenAI API Server for Qwen 3.6
After=network.target

[Service]
Type=simple
User=kran
WorkingDirectory=/home/kran/Code/kran/llm-experiments/llama.cpp

ExecStart=/home/kran/Code/kran/llm-experiments/llama.cpp/qwen3.6.bash
# Handling failures and logs
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
